# systemd/lever-runner-bot.service
[Unit]
Description=Lever-Runner Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/lever-runner
EnvironmentFile=/home/ubuntu/lever-runner/.env
Environment=PYTHONPATH=/home/ubuntu/lever-runner/src
ExecStart=/home/ubuntu/lever-runner/.venv/bin/python -m lever_runner.bot
Restart=on-failure
RestartSec=5
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
