# systemd/lever-runner-bot.service
[Unit]
Description=Lever-Runner Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/lever-runner
EnvironmentFile=/home/ubuntu/lever-runner/.env
# `lever-runner-bot` is a console-script entry point installed by
# `pip install -e .`. It calls bot:Application.run_polling().
ExecStart=/home/ubuntu/lever-runner/.venv/bin/lever-runner-bot
Restart=on-failure
RestartSec=5
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
